# systemd/battery/battery.service
[Unit]
Description=Day time wallpaper mode

[Service]
ExecStart=%h/bin/battery_notify.sh
