# systemd/battery/battery.timer
[Unit]
Description=battery notifier

[Timer]
OnUnitActiveSec=10
Persistent=true
Unit=battery.service

[Install]
WantedBy=default.target
